# 8.1/nginx/etc/supervisor.d/create_database.ini
[program:create_database]
command=/create_database.sh
user=root
autostart=true
autorestart=false
priority=10
startretry=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

